# docker/rootfs/etc/s6-overlay/s6-rc.d/qbittorrent/run
#!/command/with-contenv bash
echo "[Nexus] Starting qBittorrent..."
exec qbittorrent-nox --webui-port=8091 --profile=/config/qbittorrent
