# docker/rootfs/etc/s6-overlay/s6-rc.d/lidarr/run
#!/command/with-contenv bash
echo "[Nexus] Starting Lidarr..."
exec /opt/Lidarr/Lidarr -nobrowser -data=/config/lidarr
