# docker/rootfs/etc/s6-overlay/s6-rc.d/nexus-setup/up
#!/command/with-contenv bash

echo "[Nexus] Performing initial setup..."

# Ensure all service config directories exist
for service in qbittorrent sabnzbd sonarr radarr lidarr whisparr prowlarr overseerr flaresolverr; do
    mkdir -p "/config/$service"
done

# Ensure download subdirectories exist
mkdir -p /downloads/incomplete/torrents
mkdir -p /downloads/incomplete/usenet
mkdir -p /downloads/complete/torrents/{tv,movies,music,adult}
mkdir -p /downloads/complete/usenet/{tv,movies,music,adult}

# Ensure transcode subdirectories exist
mkdir -p /transcode/{queue,working,completed,failed,logs}

# Fix permissions for /config and /transcode
# We use PUID/PGID if provided, otherwise default to 99:100 (Unraid)
CHOWN_USER="${PUID:-99}:${PGID:-100}"
echo "[Nexus] Setting ownership to $CHOWN_USER..."
chown -R "$CHOWN_USER" /config /downloads /transcode /tv /movies /music /adult

echo "[Nexus] Setup complete."
