# docker/rootfs/etc/s6-overlay/s6-rc.d/nexus-worker/run
#!/command/with-contenv bash

echo "[Nexus] Starting Worker..."
cd /app/worker
exec /app/venv/bin/python main.py
