# docker/rootfs/etc/s6-overlay/s6-rc.d/sabnzbd/run
#!/command/with-contenv bash
echo "[Nexus] Starting SABnzbd..."
exec sabnzbd --browser 0 --server 127.0.0.1:8092 --config-file /config/sabnzbd/sabnzbd.ini
