# docker/rootfs/etc/s6-overlay/s6-rc.d/overseerr/run
#!/command/with-contenv bash
echo "[Nexus] Starting Overseerr..."
export CONFIG_DIRECTORY=/config/overseerr
exec /opt/Overseerr/overseerr
